# build/supervisor/supervisor.ini
[program:app-a]
command=npm run start
environment=PORT=8081
directory=/opt/app-a
stderr_logfile = /var/log/supervisord/app-a-stderr.log
stdout_logfile = /var/log/supervisord/app-a-stdout.log

[program:app-b]
command=npm run start
environment=PORT=8082
directory=/opt/app-b
stderr_logfile = /var/log/supervisord/app-b-stderr.log
stdout_logfile = /var/log/supervisord/app-b-stdout.log

[program:app-c]
command=npm run start
environment=PORT=8083
directory=/opt/app-c
stderr_logfile = /var/log/supervisord/app-c-stderr.log
stdout_logfile = /var/log/supervisord/app-c-stdout.log

[program:nginx]
command=/usr/sbin/nginx
pidfile=/var/run/nginx.pid
stderr_logfile = /var/log/supervisord/nginx-stderr.log
stdout_logfile = /var/log/supervisord/nginx-stdout.log
autorestart = true
